# infra/front.service
[Unit]
Description=Whats Organizer Front
After=network.target

[Service]
User=front
Group=front
WorkingDirectory=/home/front/whats-organizer-front
ExecStart=/bin/sh -lc 'export NVM_DIR="$([ -z \"${XDG_CONFIG_HOME-}\" ] && echo \"${HOME}/.nvm\" || echo \"${XDG_CONFIG_HOME}/nvm\")" && . \"$NVM_DIR/nvm.sh\" && exec node --env-file=.env build/index --host 0.0.0.0 --port 3000'
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
